# backup.service
[Unit]
Description=Run backup daily.

[Service]
Type=oneshot
ExecStart=/home/michaelstephens/git/dotfiles
